# stacks/flight-tracker/flight-tracker-poll.service
[Unit]
Description=Pull latest flight-tracker backend image and restart if changed

[Service]
Type=oneshot
WorkingDirectory=/home/colin/code/homelab
ExecStart=/bin/bash -c '\
  docker compose -f stacks/flight-tracker/compose.yaml pull --quiet backend && \
  docker compose -f stacks/flight-tracker/compose.yaml up -d --remove-orphans backend'
