# stacks/knowledge/knowledge-backup.service
[Unit]
Description=Back up knowledge Postgres database

[Service]
Type=oneshot
WorkingDirectory=/home/colin/code/homelab
Environment=KNOWLEDGE_BACKUP_DIR=/home/colin/backups/knowledge
Environment=KNOWLEDGE_BACKUP_RETENTION_DAYS=14
ExecStart=/home/colin/code/homelab/scripts/backup-knowledge-db.sh
